# etc/profile.d/FireMotD.sh
#!/usr/bin/env bash
set -v

/usr/local/bin/FireMotD --Theme Invader